# system_files/usr/lib/systemd/system/upower.service.d/10-cloudws-wsl2.conf
# upower probes /sys/class/power_supply and DBus hardware interfaces that
# don't exist in WSL2, causing an immediate SIGABRT. Skip it entirely.
# WSLInterop path check is the project-wide standard (see 0cdcf9a) because
# systemd's ConditionVirtualization=wsl detection is unreliable in containers.
[Unit]
ConditionPathExists=!/proc/sys/fs/binfmt_misc/WSLInterop
