# system_files/usr/lib/systemd/system/cloudws-nvidia-cdi.service.d/10-virt-gate.conf
[Unit]
# CloudWS: Skip in containers/WSL2
ConditionVirtualization=!container
ConditionVirtualization=!wsl
