# system_files/etc/systemd/system/dbus-broker.service.d/wsl2-fix.conf
[Service]
# Raise OOM priority so a memory-stressed WSL2 guest doesn't kill D-Bus.
# ExecStart / --audit override lives in
#   /usr/lib/systemd/system/dbus-broker.service.d/10-cloudws-no-audit.conf
# which already removes --audit for all environments.
# The previous ConditionPathExists=|/proc/sys/fs/binfmt_misc/WSLInterop
# that was here caused dbus-broker to be SKIPPED on bare metal.
OOMScoreAdjust=-500
