# system_files/etc/systemd/system/dbus-broker.service.d/wsl2-fix.conf
[Unit]
ConditionPathExists=|/proc/sys/fs/binfmt_misc/WSLInterop

[Service]
ExecStart=
ExecStart=/usr/bin/dbus-broker-launch --scope system
OOMScoreAdjust=-500
