# system_files/etc/systemd/system/dbus-broker.service.d/wsl2-fix.conf
# CloudWS — dbus-broker WSL2 workaround
# Microsoft WSL2 kernel lacks audit socket, causing dbus-broker to fail
# with exit code 1, cascading to kill ALL D-Bus dependent services.
# This drop-in strips --audit and only activates inside WSL2.

[Unit]
ConditionPathExists=|/proc/sys/fs/binfmt_misc/WSLInterop

[Service]
ExecStart=
ExecStart=/usr/bin/dbus-broker-launch --scope system
OOMScoreAdjust=-500