# system_files/usr/lib/systemd/system/dbus-broker.service.d/10-cloudws-wsl2.conf
# Skip dbus-broker on WSL2 (WSL kernels strip the audit subsystem, causing broker to crash)
[Unit]
ConditionVirtualization=!wsl
