# system_files/etc/systemd/system/auditd.service.d/10-cloudws-wsl2.conf
[Unit]
# Skip auditd in WSL2 — the kernel audit subsystem is absent there and
# auditd exits NOPERMISSION every time, exhausting its restart budget.
ConditionVirtualization=!wsl
