# usr/lib/systemd/system/systemd-networkd-wait-online.service.d/10-mios-wsl2.conf
# Skip network wait on WSL2 to prevent user session timeouts (WSL v0.1.1+)
[Unit]
ConditionVirtualization=!wsl
