# system_files/usr/lib/systemd/system/systemd-networkd-wait-online.service.d/10-mios-wsl2.conf
# Skip network wait on WSL2 to prevent user session timeouts (WSL v2.1.0+)
[Unit]
ConditionVirtualization=!wsl
