# system_files/etc/systemd/system/systemd-machined.service.d/wsl2-optional.conf
# CloudWS — systemd-machined needs cgroup features WSL2 lacks
[Unit]
ConditionVirtualization=!wsl