# system_files/usr/lib/systemd/system/libvirtd.service.d/10-cloudws.conf
# Prevent startup race condition between Cockpit and libvirt on ucore bases
[Unit]
After=libvirtd.socket
