# system_files/etc/systemd/system/audit-rules.service.d/10-cloudws-wsl2.conf
[Unit]
# Skip audit-rules in WSL2 — same reason as auditd.service.d drop-in.
ConditionVirtualization=!wsl
