# system_files/etc/systemd/system/cockpit.socket.d/listen-all.conf
# CloudWS — Force cockpit to listen on ALL interfaces
# Default cockpit.socket only listens on specific addresses
[Socket]
# Clear existing ListenStream directives first
ListenStream=
# Then bind to all IPv4 and IPv6 on port 9090
ListenStream=0.0.0.0:9090
ListenStream=[::]:9090
FreeBind=yes