# system_files/etc/systemd/system/cockpit.socket.d/listen-all.conf
# CloudWS — Cockpit listens on ALL interfaces (not just localhost)
[Socket]
ListenStream=
ListenStream=0.0.0.0:9090
ListenStream=[::]:9090
FreeBind=yes