# system_files/usr/lib/systemd/system/cloudws-wsl-firstboot.service
[Unit]
Description=CloudWS WSL2 First Boot Initialization
Documentation=https://github.com/Kabuki94/CloudWS-bootc
ConditionVirtualization=wsl
ConditionPathExists=!/var/lib/cloudws/.wsl-firstboot-done
After=local-fs.target
Before=systemd-user-sessions.service multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/cloudws/wsl-firstboot

[Install]
WantedBy=multi-user.target
