# system_files/usr/lib/systemd/system/cloudws-gpu-pv-detect.service
[Unit]
Description=CloudWS-OS Hyper-V GPU-PV Guest Detection
DefaultDependencies=no
After=systemd-modules-load.service
Before=display-manager.service
ConditionVirtualization=microsoft

[Service]
Type=oneshot
ExecStart=/usr/libexec/cloudws/gpu-pv-detect
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
