# system_files/usr/lib/systemd/system/cloudws-grd-setup.service
[Unit]
Description=CloudWS GNOME Remote Desktop Setup
ConditionVirtualization=!container
After=local-fs.target network-online.target gdm.service
Wants=network-online.target
ConditionPathExists=!/var/lib/cloudws/.grd-init-done
[Service]
Type=oneshot
ExecStart=/usr/libexec/cloudws-grd-setup
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
