# system_files/usr/lib/systemd/system/cloudws-headless.target
[Unit]
Description=CloudWS Headless role (SSH + Podman + uupd, for VMs/Podman-machine/WSL)
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes

[Install]
Alias=default.target.cloudws-headless