# system_files/usr/lib/systemd/system/cloudws-gpu-intel.service
[Unit]
Description=CloudWS Intel GPU container plumbing (i915/xe + renderD*)
ConditionVirtualization=!container
Documentation=https://github.com/Kabuki94/CloudWS-bootc/blob/main/docs/gpu-passthrough.md
ConditionPathExists=/dev/dri/renderD128
After=systemd-modules-load.service systemd-udev-settle.service
Wants=systemd-udev-settle.service
# NOTE: runs even alongside NVIDIA (hybrid/iGPU+dGPU systems are common, e.g.
# AMD 9950X3D iGPU + RTX 4090 dGPU — renderD128 is typically the AMD/Intel
# iGPU and renderD129 is the NVIDIA dGPU). Intel-specific container
# workloads (QuickSync, OpenCL via NEO) reference renderD128 explicitly.

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-/usr/sbin/modprobe -a i915 xe
ExecStartPre=/usr/bin/install -d -m 0755 /var/run/cdi
ExecStart=/usr/bin/bash -c ' \
  # Defensive perms on renderD128 only — do not touch card* nodes (video group \
  # handled by 99-cloudws-gpu.rules at udev-add time). \
  if [ -e /dev/dri/renderD128 ]; then \
    chgrp render /dev/dri/renderD128 2>/dev/null || true; \
    chmod 0660   /dev/dri/renderD128 2>/dev/null || true; \
  fi; \
  exit 0'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
