# system_files/usr/lib/systemd/system/cloudws-verify-root.service
[Unit]
Description=CloudWS Root Filesystem Verification
DefaultDependencies=no
Before=local-fs.target
After=ostree-remount.service
ConditionPathExists=/run/ostree-booted

[Service]
Type=oneshot
ExecStart=/usr/libexec/cloudws/verify-root.sh
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target