# system_files/usr/lib/systemd/system/cloudws-sriov-init.service
[Unit]
Description=CloudWS-OS Universal SR-IOV Initialization
After=network-pre.target
Wants=network-pre.target
ConditionVirtualization=!wsl
ConditionPathExists=/sys/class/net

[Service]
Type=oneshot
ExecStart=/usr/libexec/cloudws/cloudws-sriov-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
