# system_files/usr/lib/systemd/system/cloudws-podman-gc.service
[Unit]
Description=CloudWS Podman Garbage Collection
[Service]
Type=oneshot
ExecStart=/usr/bin/podman system prune -a -f
