# system_files/usr/lib/systemd/system/cloudws-ipa-enroll.service
[Unit]
Description=CloudWS FreeIPA first-boot enrollment
Documentation=file:///etc/cloudws/ipa.conf.example
Wants=network-online.target nss-lookup.target
After=network-online.target nss-lookup.target
ConditionPathExists=/etc/cloudws/ipa.conf
ConditionPathExists=!/var/lib/cloudws/ipa-enrolled

[Service]
Type=oneshot
ExecStart=/usr/libexec/cloudws/ipa-enroll
RemainAfterExit=yes
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
