# system_files/usr/lib/systemd/system/mios-mcp.service
[Unit]
ConditionVirtualization=!container
Description=Claude OS Persistent Memory MCP Server (Port 8051)
After=network.target redis.service
Wants=redis.service
# Run constantly in the background to avoid MCP "cold start" latency
# and keep SQLite databases + Vector embeddings warm in RAM.

[Service]
Type=simple
User=mios
Group=mios
# Execution wrapper that verifies SQLite vaults exist before starting the MCP listener
ExecStartPre=/usr/libexec/mios/mcp-init.sh
ExecStart=/usr/libexec/mios/mcp-server-runner
Restart=always
RestartSec=5
# Security hardening: isolate the MCP server from writing to the core OS
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/mios/mcp /var/log/mios/mcp

[Install]
WantedBy=multi-user.target
