# system_files/usr/share/containers/systemd/cloudws-guacamole.container
[Unit]
Description=CloudWS Guacamole Web App
ConditionVirtualization=!wsl
# Ensure guacd is running before the webapp starts
Requires=guacd.service
After=guacd.service

[Container]
Image=docker.io/guacamole/guacamole:latest
ContainerName=cloudws-guacamole
Network=guacamole.network
# Port 8080 matches the firewall rule in 25-firewall-ports.sh
PublishPort=8080:8080
Environment=GUACD_HOSTNAME=guacd
Environment=GUACD_PORT=4822
AutoUpdate=registry

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
