# system_files/usr/share/containers/systemd/cloudws-pxe-hub.container
[Unit]
Description=CloudWS PXE Hub (netboot.xyz)
After=network-online.target
Wants=network-online.target

[Container]
HttpProxy=false
GlobalArgs=--storage-opt=additionalimagestore=/usr/lib/bootc/storage
Image=ghcr.io/netbootxyz/netbootxyz:latest
ContainerName=cloudws-pxe-hub
# Map persistent storage for ISOs and custom menus
Volume=/var/lib/cloudws/pxe/config:/config:Z
Volume=/var/lib/cloudws/pxe/assets:/assets:Z
# Expose required PXE/Network Boot ports
# 69/udp (TFTP), 80/tcp (HTTP Assets), 8080/tcp (Web UI)
PublishPort=69:69/udp
PublishPort=80:80/tcp
PublishPort=8080:8080/tcp
# Enable host networking mode is often preferred for PXE/DHCP proxy scenarios,
# but we'll try bridge with explicit ports first for better isolation.
# If DHCP proxy is needed, users can toggle to Network=host.
Environment=PORT=80
Environment=WEB_PORT=8080
Restart=always

[Service]
# Ensure directories exist
ExecStartPre=/usr/bin/mkdir -p /var/lib/cloudws/pxe/config /var/lib/cloudws/pxe/assets

[Install]
WantedBy=multi-user.target
