# system_files/usr/share/containers/systemd/mios-nut.container
[Unit]
Description=Network UPS Tools (NUT) Distrobox
After=network-online.target
Wants=network-online.target
# Ensure NUT stops LAST during shutdown to prevent premature power-off
Before=podman.service libvirtd.service virtqemud.service virtnetworkd.service

[Container]
HttpProxy=false
Image=docker.io/library/fedora:40
ContainerName=mios-nut
# Decouple mutable NUT config from immutable host
Volume=/var/lib/mios/nut:/etc/nut:Z
Volume=/dev/bus/usb:/dev/bus/usb
# Map sockets to host for client interaction if needed
Volume=/run/nut:/run/nut:z
# Ensure hardware access
AddDevice=/dev/bus/usb
# Enable systemd inside container for internal daemons (upsd, upsmon)
Systemd=always
# Run as root to manage hardware/shutdowns
User=root
# Auto-restart on failure
Restart=always

[Service]
# Ensure the config dir exists on host before start
ExecStartPre=/usr/bin/mkdir -p /var/lib/mios/nut /run/nut

[Install]
WantedBy=multi-user.target
