# overlay/usr/share/xdg/autostart/localsearch-writeback-3.desktop
[Desktop Entry]
Type=Application
Name=localsearch-writeback-3 Override
Exec=/usr/bin/true
Hidden=true
X-GNOME-Autostart-enabled=false
